# scripts/packaging/deb/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable backend-api || true

exit 0